# .config/profile.d/_flux
#!/bin/sh
[ $(command -v flux) ] && source <(flux completion ${SHELL#\/bin\/})
