# .config/profile.d/_gh
#!/bin/sh
[ $(command -v gh) ] && source <(gh completion -s $(echo $SHELL | sed 's#^.*/##'))
