# .config/profile.d/_gh
[ $(command -v gh) ] && source <(gh completion -s ${SHELL#\/bin\/})
