# .config/profile.d/_helm
[ $(command -v helm) ] && source <(helm completion ${SHELL#\/bin\/})
