# .config/profile.d/_kubectl
[ $(command -v kubectl) ] && source <(kubectl completion ${SHELL#\/bin\/})
