# .config/profile.d/_kubectl
#!/bin/sh
[ $(command -v kubectl) ] && source <(kubectl completion $(echo $SHELL | sed 's#^.*/##'))
