# .config/profile.d/_psql
#!/bin/sh
export PATH="/opt/homebrew/opt/libpq/bin:$PATH"
