# .config/profile.d/_velero
#!/bin/sh
[ $(command -v velero) ] && velero completion zsh > "${fpath[1]}/_velero"
