# .config/profile.d/_wasmedge
#!/bin/sh
export PATH="$HOME/.wasmedge/bin:$PATH"
