# .config/profile.d/_wasmedge
export PATH="$HOME/.wasmedge/bin:$PATH"
